# confs/supervisor/config.d/frp.ini
[program:frp]
directory=/home/fagangjin/work/UbuntuScripts/bins
command=/home/fagangjin/work/UbuntuScripts/bins/frpc -c /home/fagangjin/work/UbuntuScripts/bins/frpc.ini.c8s
autostart=true
autorestart=true
stdout_logfile=/home/fagangjin/work/UbuntuScripts/bin
stderr_logfile=/home/fagangjin/work/UbuntuScripts/bin
